# stan/stan
### BEGIN INIT INFO
# Provides:          stan
# Required-Start:    
# Required-Stop:     
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start stan at boot time.
# Description:       Start stan at boot time.
### END INIT INFO

#!/bin/bash

pid=/var/run/stan.pid
. /lib/lsb/init-functions
case "$1" in
        start)
                log_begin_msg "Starting stan service"
                start-stop-daemon --start --background --exec "/usr/bin/java" -m --pidfile "$pid" -d "/home/brando/STAN/stan" -- -Xms256m -Xmx8000m -cp .:"stan.jar" it.disco.unimib.stan.main.StartWebApplication
                status=$?
		log_end_msg $status
                ;;
        stop)
                log_begin_msg "Stopping stan service"
                start-stop-daemon --oknodo --stop --pidfile "$pid"
                status=$?	
		log_end_msg $?
		rm -f $pid
                ;;
        *)
                log_success_msg "Usage: service stan start|stop"
        		exit 1
        		;;
esac
exit $?
